# fetchit/systemd/edge-autoupdate.service
[Unit]
Description=Edge Auto-Update Configuration
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c 'echo "Configuring edge-optimized update schedule..." && systemctl stop transactional-update.timer && systemctl disable transactional-update.timer && echo "Disabled daily updates for edge stability"'
ExecStart=/usr/bin/bash -c 'echo "Creating weekly update schedule for edge devices..." && cat > /etc/systemd/system/edge-update.timer << EOF
[Unit]
Description=Weekly Edge System Updates
Requires=edge-update.service

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
EOF'
ExecStart=/usr/bin/bash -c 'echo "Creating edge update service..." && cat > /etc/systemd/system/edge-update.service << EOF
[Unit]
Description=Edge System Update Service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/transactional-update --non-interactive pkg update
ExecStart=/usr/bin/systemctl reboot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
EOF'
ExecStart=/usr/bin/systemctl daemon-reload
ExecStart=/usr/bin/systemctl enable edge-update.timer
ExecStart=/usr/bin/systemctl start edge-update.timer
ExecStart=/usr/bin/echo "Edge update schedule configured: Weekly updates with 1-hour randomization"

[Install]
WantedBy=multi-user.target
