# fetchit/systemd/fetchit-restart.service
[Unit]
Description=FetchIt Auto-Restart Service
After=network.target
Wants=fetchit.service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c 'echo "Checking if FetchIt needs restart..." && if [ -f /opt/fetchit/restart-flag ]; then echo "Restarting FetchIt..." && systemctl restart fetchit.service && rm -f /opt/fetchit/restart-flag && echo "FetchIt restarted successfully"; else echo "No restart needed"; fi'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
