# fetchit/systemd/monitor.service
[Unit]
Description=FetchIt Container Monitor
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/podman ps --format json
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
